# deployment/hl-backend.service
[Unit]
Description=智学 OS 后端服务
After=network.target

[Service]
Type=simple
User=nobody
WorkingDirectory=/opt/hl-os/backend
Environment="NODE_ENV=production"
Environment="PORT=3000"
EnvironmentFile=/opt/hl-os/.env
ExecStart=/usr/bin/node dist/index.js
Restart=on-failure
RestartSec=10s
StandardOutput=journal
StandardError=journal
SyslogIdentifier=hl-backend

# 资源限制
LimitNOFILE=65536
MemoryLimit=512M

[Install]
WantedBy=multi-user.target
